# docker/supervisord.conf
; =============================================================================
; Supervisord Configuration for MCP Server
; Manages 5 processes: 4 MCP servers + 1 metrics server
; =============================================================================

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; =============================================================================
; SSE Transport Servers (HTTP Streaming via Server-Sent Events)
; =============================================================================

[program:general-sse]
command=/app/.venv/bin/python -m src.mcp_servers.mcp_general.http_main
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HTTP_PORT="8000"

[program:cobol-sse]
command=/app/.venv/bin/python -m src.mcp_servers.mcp_cobol_analysis.http_main
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=HTTP_PORT="8001"

; =============================================================================
; Streamable HTTP Transport Servers
; =============================================================================

[program:general-streamable]
command=/app/.venv/bin/python -m src.mcp_servers.mcp_general.streamable_http_main
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=STREAMABLE_HTTP_PORT="8002"

[program:cobol-streamable]
command=/app/.venv/bin/python -m src.mcp_servers.mcp_cobol_analysis.streamable_http_main
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=STREAMABLE_HTTP_PORT="8003"

; =============================================================================
; Metrics Server (Prometheus + Health Check)
; =============================================================================

[program:metrics]
command=/app/.venv/bin/python /app/scripts/start_metrics_server.py
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; =============================================================================
; Process Group (for collective management)
; =============================================================================

[group:mcp-servers]
programs=general-sse,cobol-sse,general-streamable,cobol-streamable,metrics
priority=999
